# Efficient AWS Infrastructure Provisioning with Pulumi - 01/deploy/nodejs-app.service
[Unit]
Description=Node.js Application
After=mysql-check.service
Requires=mysql-check.service

[Service]
Type=simple
User=nodejs
WorkingDirectory=/opt/app
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target